# tokyotyrant/files/etc/cron.d/tt_1978_info_pcpu
*/1 * * * * root /usr/local/sbin/tt_1978_info_pcpu.sh >>/var/log/tt_info/tt_1978_info_pcpu.log 2>&1 ;
