# tokyotyrant/files/etc/cron.d/tt_1978_info_threads
*/1 * * * * root /usr/local/sbin/tt_1978_info_threads.sh >>/var/log/tt_info/tt_1978_info_threads.log 2>&1 ;
